# Services/p1datasync.timer
[Timer]
OnCalendar=daily
Persistent=true